# systemd/raspise.service
[Unit]
Description=RaspISE — Network Access Control Server
Documentation=https://github.com/ssnrshnn/raspISE
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=60s
StartLimitBurst=5

[Service]
Type=simple
User=raspise
Group=raspise

WorkingDirectory=/opt/raspise

# Runtime environment
Environment=RASPISE_CONFIG=/etc/raspise/config.yaml
Environment=PYTHONPATH=/opt/raspise
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONDONTWRITEBYTECODE=1

# Main process
ExecStart=/opt/raspise/venv/bin/python -m raspise.main

# Graceful reload (re-reads config)
ExecReload=/bin/kill -HUP $MAINPID

# Restart policy
Restart=on-failure
RestartSec=5s

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=raspise

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/raspise /var/log/raspise /tmp
# Allow raw sockets for RADIUS/TACACS+ UDP/TCP (need CAP_NET_BIND_SERVICE for ports <1024)
AmbientCapabilities=CAP_NET_BIND_SERVICE CAP_NET_RAW
CapabilityBoundingSet=CAP_NET_BIND_SERVICE CAP_NET_RAW

# Resource limits (Pi Zero 2W has 512 MB RAM)
MemoryMax=256M
TasksMax=64

[Install]
WantedBy=multi-user.target
